# crontab
# Let's Encrypt 證書自動更新
# 每天凌晨 2 點檢查並更新證書
0 2 * * * cd /path/to/luftdeck-backend && ./renew-cert.sh >> /var/log/certbot-renew.log 2>&1
